# data/config/openwb.cron
@reboot openwb /var/www/html/openWB/runs/atreboot.sh &
