# data/config/openwb.cron
# openwb-version:1
@reboot openwb /var/www/html/openWB/runs/atreboot.sh &
