# data/config/openwb.cron
@reboot root /var/www/html/openWB/runs/atreboot.sh &
